# nginx-allow.cron
 * * * * /domain-resolver.sh /hostdns.list >/dev/null 2>&1
